# supervisor/supervisor-meritd.conf
[program:app-meritd]
command = /usr/local/bin/meritd --deamon